# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:llama]
command=/usr/local/bin/entrypoint.sh llama
autostart=true
autorestart=true
startsecs=10
startretries=1000000
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:tts]
command=python -m uvicorn tts.server:app --host 127.0.0.1 --port 8082 --no-access-log
directory=/app
autostart=true
autorestart=true
startsecs=10
startretries=1000000
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:gateway]
command=python -m uvicorn gateway.app:app --host 0.0.0.0 --port 8000 --no-access-log
directory=/app
autostart=true
autorestart=true
startsecs=5
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
