# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:llama]
command=/usr/local/bin/entrypoint.sh llama
autostart=true
autorestart=true
startsecs=10
startretries=1000000
priority=10
redirect_stderr=true
stdout_logfile=/models/llama.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2

[program:gateway]
command=python -m uvicorn gateway.app:app --host 0.0.0.0 --port 8000 --no-access-log
directory=/opt/app
autostart=true
autorestart=true
startsecs=5
priority=30
redirect_stderr=true
stdout_logfile=/models/api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
